# web/crontab
# choose a random minute for renewing certificate.
# certbot will be thankful for it.
39 0 * * * certbot-auto renew --no-self-upgrade >/dev/null 2>&1
